# scripts/ai.tpsdev.harper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.tpsdev.harper</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>FLAIR_DIR/node_modules/harperdb/bin/harper.js</string>
    <string>run</string>
    <string>FLAIR_DIR</string>
  </array>

  <key>WorkingDirectory</key>
  <string>FLAIR_DIR</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>HOME_DIR/.tps/logs/harper.log</string>

  <key>StandardErrorPath</key>
  <string>HOME_DIR/.tps/logs/harper.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>HOME_DIR</string>
    <key>ROOTPATH</key>
    <string>HOME_DIR/.harper/flair</string>
    <key>HTTP_PORT</key>
    <string>9926</string>
    <key>OPERATIONSAPI_NETWORK_PORT</key>
    <string>9925</string>
  </dict>
</dict>
</plist>
