# scripts/ai.tpsdev.flair-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.tpsdev.flair-watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/squeued/ops/flair/scripts/harper-watchdog.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/squeued/.tps/logs/harper-watchdog.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/squeued/.tps/logs/harper-watchdog.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/squeued</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HARPER_PORT</key>
    <string>9926</string>
    <key>LAUNCHD_LABEL</key>
    <string>ai.tpsdev.flair</string>
  </dict>
</dict>
</plist>
