# scripts/ai.tpsdev.flair.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.tpsdev.flair</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>FLAIR_DIR/node_modules/harperdb/bin/harper.js</string>
    <string>dev</string>
    <string>FLAIR_DIR</string>
  </array>

  <key>WorkingDirectory</key>
  <string>FLAIR_DIR</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>HOME_DIR/.tps/logs/harper.log</string>

  <key>StandardErrorPath</key>
  <string>HOME_DIR/.tps/logs/harper.error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>HOME_DIR</string>
    <key>HARPER_SET_CONFIG</key>
    <string>{"rootPath":"HOME_DIR/.harper/flair","http":{"port":9926},"operationsApi":{"network":{"port":9925}},"authentication":{"operationsAdminPassword":"ADMIN_TOKEN_PLACEHOLDER"}}</string>

  <key>FLAIR_ADMIN_TOKEN</key>
  <string>ADMIN_TOKEN_PLACEHOLDER</string>
  </dict>
</dict>
</plist>
